# embedding-service/embedding-service.service
[Unit]
Description=Magento Vector Search - Embedding Service (multilingual-e5-base)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/app/code/Kkkonrad/VectorSearch/embedding-service
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=PORT=3000
Environment=HOST=127.0.0.1
Environment=NODE_ENV=production
Environment=INTRA_OP_THREADS=4
Environment=INTER_OP_THREADS=4
EnvironmentFile=-/etc/default/magento-vector-search
# Allow enough RAM for the ONNX model (~500MB)
MemoryMax=1G
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/html/app/code/Kkkonrad/VectorSearch/embedding-service/models
StandardOutput=journal
StandardError=journal
SyslogIdentifier=embedding-service

[Install]
WantedBy=multi-user.target
